# root/etc/s6-overlay/s6-rc.d/10-homers/run
#!/command/execlineb -P
with-contenv
s6-setuidgid homers
homers --config /app/config.toml
